# 2026-04-26/role-and-objective-deploy-and-harden/ops/systemd/trading-engine-paper.service
[Unit]
Description=Trading Engine Paper Strategy Service
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
User=trader
Group=trader
WorkingDirectory=/opt/trading-system/app
EnvironmentFile=/opt/trading-system/shared/.env.paper
ExecStart=/opt/trading-system/app/.venv/bin/python -m trading_system.trading.paper_strategy_runner --env-file /opt/trading-system/shared/.env.paper --interval-seconds 900
Restart=on-failure
RestartSec=30
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/trading-system/shared /opt/trading-system/app
LimitNOFILE=65535
StandardOutput=append:/opt/trading-system/shared/logs/paper.out.log
StandardError=append:/opt/trading-system/shared/logs/paper.err.log

[Install]
WantedBy=multi-user.target
